# displays/systemd/autorandr.service.d/override.conf
[Service]
# Override the stock fallback: the shipped unit uses `--default default`, which
# points at a profile named "default" that doesn't exist — so a hotplugged
# monitor with no saved profile stayed blank. `horizontal` is a built-in virtual
# profile that extends whatever is connected left-to-right, so ANY new monitor
# now lights up automatically. Installed by ~/setup/displays/install.sh.
ExecStart=
ExecStart=/usr/bin/autorandr --batch --change --default horizontal
