# displays/systemd/autorandr-lid-listener.service.d/override.conf
[Service]
# Same fix as autorandr.service, but for lid events: the shipped lid-listener
# also falls back to `--default default` (nonexistent), so closing the lid did
# nothing. `horizontal` makes autorandr re-evaluate; because autorandr marks the
# internal panel (eDP) as disconnected while the lid is closed AND an external is
# connected, applying `horizontal` turns eDP OFF and keeps the externals. Opening
# the lid fires again and brings eDP back. Installed by displays/install.sh.
#
# NOTE: this ExecStart must stay byte-for-byte identical to the shipped unit
# except for the trailing `--default default` -> `--default horizontal`.
ExecStart=
ExecStart=sh -c "stdbuf -oL libinput debug-events | grep -E --line-buffered '^[[:space:]-]+event[0-9]+[[:space:]]+SWITCH_TOGGLE[[:space:]]' | while read line; do autorandr --batch --change --default horizontal; done"
